# scripts/templates/systemd/vpn-store-web.service.template
[Unit]
Description=VPN Store web service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_GROUP}}
WorkingDirectory={{INSTALL_DIR}}
EnvironmentFile={{INSTALL_DIR}}/.env
ExecStart={{INSTALL_DIR}}/venv/bin/gunicorn core.wsgi:application --bind 127.0.0.1:{{GUNICORN_PORT}}
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
KillSignal=SIGQUIT
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
