# scripts/templates/systemd/timers/vpn-store-panel-daily-usage.timer.template
[Unit]
Description=Calculate VPN Store panel daily usage

[Timer]
OnCalendar=*-*-* 00:20:00
Persistent=true

[Install]
WantedBy=timers.target
