# scripts/templates/systemd/timers/vpn-store-panel-usage-snapshots.timer.template
[Unit]
Description=Collect VPN Store panel usage snapshots periodically

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target
